# debian/koppeln.service
[Unit]
Description=A minimalistic dynamic DNS server.

[Service]
ExecStart=/usr/bin/koppeln

[Install]
WantedBy=multi-user.target
